# config/systemd/kos-zram.service
[Unit]
Description=KlipperOS-AI zram swap with zstd compression
After=local-fs.target
Before=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/klipperos-ai/scripts/setup-zram.sh
ExecStop=/usr/sbin/swapoff /dev/zram0
ExecStop=/bin/bash -c 'echo 1 > /sys/block/zram0/reset'

[Install]
WantedBy=sysinit.target
